# systemd/system/snap-core-13250.mount
[Unit]
Description=Mount unit for core, revision 13250
Before=snapd.service
After=zfs-mount.service

[Mount]
What=/var/lib/snapd/snaps/core_13250.snap
Where=/snap/core/13250
Type=squashfs
Options=nodev,ro,x-gdu.hide,x-gvfs-hide
LazyUnmount=yes

[Install]
WantedBy=multi-user.target
